# Docker/CronTab.cron
0 * * * * /usr/sbin/python /app/code/main.py
* * * * * echo $(date) > /app/now.txt
